# resources/services/ngafid-event-observer.service
[Unit]
Description=NGAFID Event Observer
After=network-online.target kafka.service
PartOf=ngafid.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/opt/ngafid/ngafid2.0
ExecStart=/bin/bash -ac '. /etc/ngafid.conf; exec run/prod/kafka/event_observer'

[Install]
WantedBy=ngafid.target
